# fs-root/etc/services.d/ss/run
#!/command/with-contenv bash
#
#
/usr/bin/docker_ss.sh

