# fs-root/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/execlineb -P
#

/usr/bin/docker_sshd.sh
