# fs-root/etc/s6-overlay/s6-rc.d/microsocks/run
#!/command/execlineb -P
#
/usr/bin/docker_microsocks.sh

